# disk-health/disk-health.timer
[Unit]
Description=Run disk-health periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
