# disk-health/smart-selftest-short.service
[Unit]
Description=Run weekly SMART short tests

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/smart-selftest-short.sh